# nginx-stream-binary/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn -b unix:app.sock -w 4 --chdir /app app:app

[program:nginx]
command=nginx -g "daemon off;"