# spark-cluster/base/ubuntu/etc/supervisor/conf.d/ssh.conf
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stderr_logfile=/var/log/sshd.err.log
stdout_logfile=/var/log/sshd.out.log