# sphinx-autobuild/ubuntu/etc/supervisor/conf.d/autobuild.conf
[program:autobuild]
command=/usr/local/sbin/sphinx-autobuild.sh
priority=1
startsecs=0
exitcodes=0