# java-jupyter/ubuntu/etc/supervisor/conf.d/all.conf
[program:jupyter]
command=/usr/local/conda/bin/jupyter lab --ip='*' --NotebookApp.token='' --NotebookApp.password='' --no-browser --allow-root
priority=999